# BuildSupport/Generated/CoreMonitorHelper-Launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>AssociatedBundleIdentifiers</key>
    <array>
        <string>CoreTools.Core-Monitor-CLI.BlessHost</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>Label</key>
    <string>ventaphobia.smc-helper</string>
    <key>MachServices</key>
    <dict>
        <key>ventaphobia.smc-helper</key>
        <true/>
    </dict>
    <key>Program</key>
    <string>/Library/PrivilegedHelperTools/ventaphobia.smc-helper</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/PrivilegedHelperTools/ventaphobia.smc-helper</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
